# docker/services/app/supervisord.conf
[supervisord]
nodaemon=true

[program:dev-web-server]
command=/bin/bash -c "cd /project/app && ./yii serve/index 0.0.0.0 -p 8000 -t backend/web --interactive=0"